# assignment_04/45-kamashis-usb.rules
SUBSYSTEM=="input", SUBSYSTEMS=="usb", ACTION=="add", RUN+="/sbin/insmod /root/home/Little-penguin/ass04/main.ko"
